# udev/64-btrfs-scan.rules
SUBSYSTEM=="block" RUN+="/sbin/btrfsctl -A %N"
